# borgbackup.timer
[Unit]
Description=Run borgbackup daily at eight

[Timer]
OnCalendar=22:00
 
[Install]
WantedBy=timers.target
